# usda_api/systemd_usda_api.service.example
[Unit]
Description=USDA API (Node.js)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/usda_api
EnvironmentFile=/opt/usda_api/.env
Environment=KEYSTORE_PATH=/var/lib/usda_api/keys.json
Environment=ADMIN_STORE_PATH=/var/lib/usda_api/admin_users.json
StateDirectory=usda_api
StateDirectoryMode=0700
ExecStart=/usr/bin/env node /opt/usda_api/src/server.js
Restart=on-failure
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/usda_api

[Install]
WantedBy=multi-user.target
